# scripts/oci/oci-a1-capture.service
[Unit]
Description=OCI A1.Flex Capture Retry (task-board #206)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Override these in /etc/systemd/system/oci-a1-capture.service.d/override.conf
Environment=A1_REGION=us-sanjose-1
Environment=A1_AD=1
Environment=A1_SHAPE=VM.Standard.A1.Flex
Environment=A1_OCPUS=2
Environment=A1_RAM_GB=12
Environment=A1_INSTANCE_NAME=gateway-a1
Environment=A1_MAX_RETRIES=200
Environment=A1_BACKOFF_MIN_S=60
Environment=A1_BACKOFF_MAX_S=3600
Environment=A1_LOG_FILE=/var/log/oci-a1-capture.log
# Required (no defaults — fill these in or set via EnvironmentFile):
# Environment=A1_IMAGE_OCID=ocid1.image.oc1...
# Environment=A1_SUBNET_OCID=ocid1.subnet.oc1...
# Environment=A1_COMPARTMENT_OCID=ocid1.tenancy.oc1...
# Environment=A1_SSH_KEY_OCID=ocid1...
EnvironmentFile=-/etc/oci-a1-capture.conf

ExecStart=/home/osees/CodingProjects/scripts/oci/a1-capture-retry.sh

# Restart policy: keep trying until success, with 5min floor between services
Restart=on-failure
RestartSec=300

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log /tmp
PrivateTmp=true
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
