# scripts/reorder-chains.service
[Unit]
Description=LLM-API-Key-Proxy: dynamic fallback chain reorder (U-formula)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Env overrides — copy to /etc/systemd/system/ and edit as needed.
Environment=TELEMETRY_DB_PATH=/dev/shm/telemetry.db
Environment=REORDER_WINDOW_H=24
Environment=REORDER_MIN_SAMPLES=5
Environment=REORDER_MAX_TPS=3000
Environment=REORDER_MAX_TTFT_MS=30000
# U-formula inputs (benchmark scores + provider category metadata).
Environment=LLM_BENCHMARK_DB=/home/ubuntu/gateway-data/benchmark.db
Environment=LLM_PROVIDERS_DB=/home/ubuntu/gateway-data/providers.db
# Path to repo + config
Environment=REPO_PATH=/home/ubuntu/LLM-API-Key-Proxy
# ponytail: systemd does NOT expand ${REPO_PATH} in WorkingDirectory — must hardcode
WorkingDirectory=/home/ubuntu/LLM-API-Key-Proxy
ExecStart=/usr/bin/python3 ${REPO_PATH}/scripts/reorder_chains.py \
  --config ${REPO_PATH}/config/virtual_models.yaml \
  --telemetry-db ${TELEMETRY_DB_PATH} \
  --tier-config ${REPO_PATH}/config/tier_config.yaml \
  --use-u-formula \
  --window-h ${REORDER_WINDOW_H} \
  --min-samples ${REORDER_MIN_SAMPLES}

# Restart gateway to pick up new config (optional; remove if hot-reload is wired).
# Post-restart leaves a 5s window where the gateway is unavailable.
ExecStartPost=/bin/sh -c 'systemctl restart llm-gateway.service || true'

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/LLM-API-Key-Proxy/config
PrivateTmp=true

[Install]
WantedBy=multi-user.target
